# garmincoach/rootfs/etc/s6-overlay/s6-rc.d/garmincoach/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# SPDX-License-Identifier: MIT
# ==============================================================================
# GarminCoach — s6-overlay run script
# ==============================================================================

declare -r DATA_DIR="/data"
declare -r LOG_DIR="/app/logs"

mkdir -p "${LOG_DIR}" "${DATA_DIR}"

# ─── Read addon configuration ─────────────────────────────────────────────────
GARMIN_EMAIL=$(bashio::config 'garmin_email')
GARMIN_PASSWORD=$(bashio::config 'garmin_password')
AI_BACKEND=$(bashio::config 'ai_backend' || echo "ha_conversation")
OLLAMA_URL=$(bashio::config 'ollama_url' || echo "")
SYNC_INTERVAL=$(bashio::config 'sync_interval_minutes' || echo "60")

export GARMIN_EMAIL GARMIN_PASSWORD
export AI_BACKEND
export OLLAMA_URL
export OLLAMA_MODEL="${OLLAMA_MODEL:-gpt-oss:20b}"
export DATABASE_URL="postgresql://postgres@127.0.0.1:5432/garmincoach"
export POSTGRES_URL="postgresql://postgres@127.0.0.1:5432/garmincoach"
export DEV_BYPASS_AUTH="true"
export NODE_ENV="production"
export AUTH_SECRET="ha-addon-local-secret-$(cat /proc/sys/kernel/random/boot_id)"
export AUTH_DISCORD_ID="unused"
export AUTH_DISCORD_SECRET="unused"

# HA Supervisor API token (auto-injected by Supervisor for addons with homeassistant_api: true)
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN:-}"
export HA_BASE_URL="http://supervisor/core"

# ─── Wait for PostgreSQL and initialize database ─────────────────────────────
bashio::log.info "Waiting for PostgreSQL..."
for i in $(seq 1 30); do
    if pg_isready -h 127.0.0.1 -p 5432 -q; then
        break
    fi
    sleep 1
done

# Create database if it doesn't exist
if ! psql -h 127.0.0.1 -U postgres -lqt | cut -d \| -f 1 | grep -qw garmincoach; then
    bashio::log.info "Creating garmincoach database..."
    createdb -h 127.0.0.1 -U postgres garmincoach
fi

# Push schema to database (creates tables if missing)
bashio::log.info "Pushing database schema..."
cd /app/packages/db && POSTGRES_URL="${POSTGRES_URL}" npx drizzle-kit push --force 2>&1 || bashio::log.warning "Schema push had warnings"

# ─── Data migration: fix userId for synced data ──────────────────────────────
psql -h 127.0.0.1 -U postgres -d garmincoach -c \
  "UPDATE daily_metric SET user_id = 'seed-user-001' WHERE user_id = 'addon-user';" 2>/dev/null
psql -h 127.0.0.1 -U postgres -d garmincoach -c \
  "UPDATE activity SET user_id = 'seed-user-001' WHERE user_id = 'addon-user';" 2>/dev/null

# ─── Start Garmin sync in background ──────────────────────────────────────────
TOKEN_DIR="/data/garmin-tokens"
# Remove initial sync marker to force 180-day backfill after userId fix
rm -f "${TOKEN_DIR}/.initial_sync_done"
if [ -n "${GARMIN_EMAIL}" ] && [ -n "${GARMIN_PASSWORD}" ]; then
    bashio::log.info "Starting Garmin sync (every ${SYNC_INTERVAL} minutes)..."
    (
        while true; do
            bashio::log.info "Syncing Garmin data..."
            python3 /app/scripts/garmin-sync.py 2>&1 || bashio::log.warning "Garmin sync failed"
            sleep "$((SYNC_INTERVAL * 60))"
        done
    ) &
elif [ -f "${TOKEN_DIR}/oauth1_token.json" ] && [ -f "${TOKEN_DIR}/oauth2_token.json" ]; then
    bashio::log.info "Found saved Garmin tokens in ${TOKEN_DIR}"
    bashio::log.info "Starting Garmin sync from saved tokens (every ${SYNC_INTERVAL} minutes)..."
    (
        while true; do
            bashio::log.info "Syncing Garmin data..."
            python3 /app/scripts/garmin-sync.py 2>&1 || bashio::log.warning "Garmin sync failed"
            sleep "$((SYNC_INTERVAL * 60))"
        done
    ) &
else
    bashio::log.warning "No Garmin credentials or saved tokens — skipping auto-sync"
    bashio::log.info "Use Settings page to log in, or run scripts/generate-garmin-tokens.py locally"
fi

# ─── Configure AI backend ─────────────────────────────────────────────────────
case "${AI_BACKEND}" in
    ha_conversation)
        if [ -n "${SUPERVISOR_TOKEN}" ]; then
            bashio::log.info "AI backend: HA Conversation API (OpenClaw/Claude)"
            bashio::log.info "  Using Supervisor token for HA API access"
        else
            bashio::log.warning "AI backend set to ha_conversation but no Supervisor token — falling back to rules-based"
            export AI_BACKEND="none"
        fi
        ;;
    ollama)
        if [ -n "${OLLAMA_URL}" ]; then
            if curl -s -o /dev/null -w "%{http_code}" "${OLLAMA_URL}/api/tags" | grep -q "200"; then
                bashio::log.info "AI backend: Ollama at ${OLLAMA_URL}"
            else
                bashio::log.warning "Ollama not reachable at ${OLLAMA_URL} — falling back to rules-based"
                export AI_BACKEND="none"
            fi
        else
            bashio::log.warning "AI backend set to ollama but no URL configured — falling back to rules-based"
            export AI_BACKEND="none"
        fi
        ;;
    none|*)
        bashio::log.info "AI backend: rules-based (no LLM)"
        export AI_BACKEND="none"
        ;;
esac

# ─── Start Next.js ────────────────────────────────────────────────────────────
bashio::log.info "Starting GarminCoach..."
cd /app/apps/nextjs

# Start Next.js on internal port (not exposed to ingress)
export HOSTNAME="127.0.0.1"
export PORT="3001"
node server.js &
NEXT_PID=$!

# Wait for Next.js to be ready
bashio::log.info "Waiting for Next.js..."
for i in $(seq 1 30); do
    if curl -s -o /dev/null http://127.0.0.1:3001/; then
        break
    fi
    sleep 1
done

# Start ingress proxy on port 3000 (HA ingress connects here)
# Rewrites HTML asset paths to work behind /api/hassio_ingress/<token>/
bashio::log.info "Starting ingress proxy..."
export PORT="3000"
export NEXT_INTERNAL_PORT="3001"
exec node /app/ingress-proxy.js
